# docker/lilac-core.service
[Unit]
Description=Lilac Core
Requires=user@@LILAC_UID@.service
After=network.target user@@LILAC_UID@.service
ConditionPathExists=!/run/lilac/verify-only
SuccessAction=exit-force
FailureAction=exit-force

[Service]
Type=simple
User=lilac
Group=lilac
WorkingDirectory=/app
Environment=HOME=/home/lilac
Environment=PATH=/data/bin:/data/.npm-global/bin:/home/lilac/.local/bin:/home/lilac/.bun/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=XDG_RUNTIME_DIR=/run/user/@LILAC_UID@
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/@LILAC_UID@/bus
Environment=LILAC_UID=@LILAC_UID@
Environment=LILAC_USER=lilac
EnvironmentFile=/run/lilac/container.env
ExecStart=/home/lilac/.bun/bin/bun apps/core/src/runtime/main.ts
KillMode=mixed
KillSignal=SIGTERM
TasksMax=infinity
TimeoutStopSec=20s
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
